# briz/update-motd.d/0-logo
#!/bin/bash
clear
cat etc/update-motd.d/logo
